# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
stderr_logfile=/var/log/php-fpm.err.log
stdout_logfile=/var/log/php-fpm.out.log

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log

[program:laravel-startup]
command=/bin/sh -c "php /var/www/artisan migrate --force && php /var/www/artisan optimize"
autostart=true
autorestart=false
startsecs=0
stderr_logfile=/var/log/laravel-startup.err.log
stdout_logfile=/var/log/laravel-startup.out.log
